# systemd/backup-notify@.service
[Unit]
Description=Backup failure notifier for (%i)

[Service]
Type=oneshot
ExecStart=/usr/local/bin/backup-notify.sh backup %i failure

