# systemd/borgmatic-check@.service
# Override file over /etc/systemd/system/borgmatic@.service
# Has to be a templated unit (@) because the original is templated unit

[Unit]
Description=borgmatic weekly consistency checks (all configs)
OnFailure=

[Service]
ExecStart=
ExecStart=systemd-inhibit --who="borgmatic-check" --what="sleep:shutdown" --why="Prevent interrupting scheduled backup check" /usr/local/bin/borgmatic-check-wrapper.sh

ExecStartPost=
ExecStopPost=

