# systemd/borgmatic@vaultwarden.timer
[Unit]
Description=Schedule: borgmatic backup (vaultwarden)

[Timer]
OnCalendar=11:30
Persistent=true

[Install]
WantedBy=timers.target

