# systemd/borgmatic-check.timer
[Unit]
Description=Schedule: borgmatic check (all)

[Timer]
OnCalendar=Sun *-*-* 13:00:00
Persistent=true

[Install]
WantedBy=timers.target

